# 04-systemd/devops-demo.service
[Unit]
Description=DevOps Demo Service

[Service]
ExecStart=/home/yarin/devops-demo/demo.sh
Restart=always

[Install]
WantedBy=multi-user.target

